# node-exporter.service
[Unit]
Description=node_exporter
[Service]
ExecStart=/opt/node_exporter/node_exporter
[Install]
WantedBy=multi-user.target
